# assets/clippard.service
[Unit]
Description=clipperd
After=default.target

[Service]
User=root
Group=root
#Type=oneshot
#RemainAfterExit=yes
Type=forking
PIDFile=/var/lib/clipperd/clipperd.pid
Restart=no
ExecStart=/bin/bash -c "/var/lib/clipperd/clipperd"
ExecStop=kill $(cat /var/lib/clipperd/clipperd.pid) && rm /var/lib/clipperd/clipperd.pid
